# whisper.service
[Unit]
Description=Whisper Transcription Service
After=default.target

[Service]
Type=simple
ExecStart=%EXEC_START%
WorkingDirectory=%WORKING_DIR%
Restart=on-failure
RestartSec=5
EnvironmentFile=%ENV_FILE%
# Add any other necessary environment variables here

[Install]
WantedBy=default.target
